# www/fabio/pkg-plist
%%ETCDIR%%.properties.sample
etc/rc.d/fabio
sbin/fabio
