# dns/knot-resolver/kresd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: kresd
# REQUIRE: NETWORKING
# BEFORE: SERVERS
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable knot-resolver:
#
# kresd_enable="YES":	Set to YES to enable kresd.
#			Set to NO by default.
# kresd_config="":	Set to %%ETCDIR%%/kresd.conf
#			by default.
#

. /etc/rc.subr

name=kresd
rcvar=kresd_enable

load_rc_config ${name}

# set defaults

kresd_enable=${kresd_enable:-"NO"}
kresd_rundir=${kresd_rundir:-"/var/db/${name}"}
kresd_config=${kresd_config:-"%%ETCDIR%%/${name}.conf"}
user=${name}
group=${name}

pidfile=/var/run/${name}.pid

command=/usr/sbin/daemon
actual_command=%%PREFIX%%/sbin/${name}
command_args="-c -f -r -P ${pidfile} ${actual_command} -c ${kresd_config} -n ${kresd_rundir}"

required_files=${kresd_config}

run_rc_command "$1"
