# dns/dhisd/pkg-plist
sbin/dhisd
sbin/dhisgenkeys
@(,,600) etc/dhis/dhis.conf.sample
@(,,600) etc/dhis/dhis.db.sample
etc/rc.d/dhisd.sh
@dir(,,700) etc/dhis
