# databases/mongodb40/pkg-plist
bin/install_compass
bin/mongo
bin/mongod
bin/mongos
etc/rc.d/mongod
@sample etc/mongodb.conf.sample
